# core/federation-dashboard.service
[Unit]
Description=S30-PRO Federation Dashboard (v2.1.0 - Master Link)
RequiresMountsFor=/mnt/raid1
After=network-online.target seestar_env_lock.service mnt-raid1.mount
Wants=network-online.target

[Service]
User=ed
Group=ed
WorkingDirectory=/home/ed/seestar_organizer
# The Master Link: Ensuring the script sees the whole project
Environment=PYTHONPATH=/home/ed/seestar_organizer
ExecStart=/home/ed/.pyenv/versions/ssc-3.13.5/bin/python3 /home/ed/seestar_organizer/core/dashboard/dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
