# core/seeing-scraper.timer
[Unit]
Description=Run Seeing Scraper Hourly

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
