# systemd/seevar-orchestrator.service
[Unit]
Description=SeeVar Science Orchestrator
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/seevar
Environment=PYTHONPATH=%h/seevar
ExecStart=%h/seevar/.venv/bin/python3 core/flight/orchestrator.py
Restart=always
RestartSec=15

[Install]
WantedBy=default.target
