# systemd/seevar-weather.service
[Unit]
Description=SeeVar Weather Sentinel
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/seevar
Environment=PYTHONPATH=%h/seevar
ExecStart=%h/seevar/.venv/bin/python3 core/preflight/weather.py
Restart=always
RestartSec=30
