# systemd/seevar-dashboard.service
[Unit]
Description=SeeVar Dashboard (Sovereign Link)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/seevar
Environment=PYTHONPATH=%h/seevar
ExecStart=%h/seevar/.venv/bin/python3 core/dashboard/dashboard.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
